# security-hardening/systemd/pb-security-hardening-check.service
[Unit]
Description=PB - Security hardening check (daily)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
UMask=0027
Environment=SYSTEMD_UNIT=%n
StandardInput=null
TimeoutStartSec=15min
KillMode=process

SupplementaryGroups=msmtp
ExecPaths=/bin /usr/bin

ExecStart=/usr/local/libexec/pb-maintenance/security-hardening-check.sh --check

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pb-security-hardening

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
ProtectSystem=strict

ReadWritePaths=/backup/security-logs /tmp /etc/msmtp