# overrides/pblinuxutility/pb-check-for-updates.service.d/no-namespace.conf
# no-namespace.conf — pblinuxutility host override
#
# Disables all directives that require CLONE_NEWNS (mount namespace support).
# These fail with exit 226 (EXIT_NAMESPACE) on this host, or cause read-only
# filesystem errors when the namespace setup is blocked but directives are
# partially applied.
#
# Note: empty-string resets (e.g. ProtectSystem=) are not accepted by this
# host's systemd version for enumerated/boolean directives — explicit
# permissive values are used instead.
#
# This file is installed to:
#   /etc/systemd/system/pb-check-for-updates.service.d/no-namespace.conf
#
# Retained from base unit (no namespace required):
#   LockPersonality=true
#   RestrictRealtime=true
#   RestrictSUIDSGID=true
#
# See: overrides/pblinuxutility/README.md
# See: DEV-GUIDE.md §6 KFC-R02

[Service]
ProtectSystem=no
ProtectHome=no
ProtectKernelModules=false
ProtectKernelTunables=false
ReadOnlyPaths=
ReadWritePaths=
