# overrides/pblinuxutility/pb-security-hardening-check-monthly.service.d/no-namespace.conf
# no-namespace.conf — pblinuxutility host override
#
# Drops all directives that require CLONE_NEWNS (mount namespace support).
# These fail with exit 226 (EXIT_NAMESPACE) on this host.
#
# This file is installed to:
#   /etc/systemd/system/pb-security-hardening-check-monthly.service.d/no-namespace.conf
#
# Retained from base unit (no namespace required):
#   NoNewPrivileges=true
#   ProtectControlGroups=true
#   LockPersonality=true
#   RestrictNamespaces=true
#   RestrictRealtime=true
#   RestrictSUIDSGID=true
#
# See: overrides/pblinuxutility/README.md
# See: DEV-GUIDE.md §6 KFC-R02

[Service]
ProtectSystem=
ProtectHome=
PrivateTmp=
ProtectKernelModules=
ProtectKernelTunables=
ReadWritePaths=
