# overrides/pblinuxutility/pb-security-hardening-check-monthly.service.d/no-namespace.conf
# no-namespace.conf — pblinuxutility host override
#
# Drops mount-namespace and kernel-capability sandbox directives that require
# CLONE_NEWNS.  These fail with exit 226 (EXIT_NAMESPACE) on this host because
# its kernel/hypervisor does not permit unprivileged mount namespaces.
#
# This file is installed to:
#   /etc/systemd/system/pb-security-hardening-check-monthly.service.d/no-namespace.conf
#
# See: overrides/pblinuxutility/README.md
# See: DEV-GUIDE.md §6 KFC-R02

[Service]
ProtectSystem=
PrivateTmp=
ProtectKernelModules=
ProtectKernelTunables=
