# overrides/pblinuxutility/pb-security-hardening-check-monthly.service.d/no-namespace.conf
# no-namespace.conf — pblinuxutility host override
#
# Disables all directives that require CLONE_NEWNS (mount namespace support).
#
# This file is installed to:
#   /etc/systemd/system/pb-security-hardening-check-monthly.service.d/no-namespace.conf
#
# See: overrides/pblinuxutility/README.md
# See: DEV-GUIDE.md §6 KFC-R02

[Service]
ProtectSystem=no
ProtectHome=no
PrivateTmp=false
ProtectKernelModules=false
ProtectKernelTunables=false
ReadWritePaths=
